# upstart/.config/upstart/mpd-sound-menu.conf
description "mpd-sound-menu"

start on started mpd
stop on stopping mpd

respawn

exec $HOME/dev/go/bin/mpd-sound-menu
