# upstart/.config/upstart/mpd.conf
description "mpd"

start on desktop-start
stop on desktop-end

respawn
expect fork

exec mpd
