# upstart/.config/upstart/mpdas.conf
description "mpdas"

start on started mpd
stop on stopping mpd

respawn

exec mpdas -c $HOME/.mpdasrc
