# upstart/.config/upstart/mpdas-trobble.conf
description "mpdas-trobble"

start on started mpd
stop on stopping mpd

respawn

exec $HOME/dev/bin/mpdas-trobble
