# upstart/.config/upstart/mpdrobble.conf
description "mpdrobble"

start on started mpd
stop on stopping mpd

respawn

exec $HOME/dev/go/bin/mpdrobble \
       --config ~/.config/mpdrobble/config.toml
