# upstart/.config/upstart/redshift.conf
description "redshift"

start on desktop-start
stop on desktop-end

respawn

exec redshift
