# upstart/.config/upstart/xbindkeys.conf
description "xbindkeys"

start on desktop-start
stop on desktop-end

respawn
expect daemon

exec xbindkeys -f $HOME/.xbindkeysrc
